# autostart/mopidy.desktop
[Desktop Entry]
Type=Application
Name=Mopidy                          
Exec=mopidy
