# systemd/projektor.service
[Unit]
Description=Projektori server

[Service]
ExecStart=/opt/projektor/server.py

[Install]
WantedBy=multi-user.target

